# dist/aetherbus.service.in
[Unit]
Description=AetherBus Serial Port Interceptor (v@PROJECT_VERSION@)
After=network.target

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/aetherbus
Restart=on-failure
RestartSec=5s
SupplementaryGroups=dialout

[Install]
WantedBy=graphical.target
